# supervisord.conf
[supervisord]
nodaemon=true
logfile=/home/container/logs/supervisord.log
logfile_maxbytes=10MB
loglevel=info
pidfile=/tmp/supervisord.pid

[unix_http_server]
file=/tmp/supervisor.sock

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

; ---------------------------------------------------------------------------
; feed-requests
; ---------------------------------------------------------------------------
[program:feed-requests]
command=node dist/main.js
directory=/app/services/feed-requests
autostart=true
autorestart=true
startretries=10
priority=10
stdout_logfile=/home/container/logs/feed-requests.log
stdout_logfile_maxbytes=5MB
stderr_logfile=/home/container/logs/feed-requests-err.log
stderr_logfile_maxbytes=5MB
environment=NODE_ENV="production"

; ---------------------------------------------------------------------------
; user-feeds-next
; ---------------------------------------------------------------------------
[program:user-feeds]
command=node ./dist/index.js
directory=/app/services/user-feeds-next
autostart=true
autorestart=true
startretries=10
priority=20
stdout_logfile=/home/container/logs/user-feeds.log
stdout_logfile_maxbytes=5MB
stderr_logfile=/home/container/logs/user-feeds-err.log
stderr_logfile_maxbytes=5MB
environment=NODE_ENV="production"

; ---------------------------------------------------------------------------
; discord-rest-listener
; ---------------------------------------------------------------------------
[program:discord-rest-listener]
command=node build/app.js
directory=/app/services/discord-rest-listener
autostart=true
autorestart=true
startretries=10
priority=30
stdout_logfile=/home/container/logs/discord-rest-listener.log
stdout_logfile_maxbytes=5MB
stderr_logfile=/home/container/logs/discord-rest-listener-err.log
stderr_logfile_maxbytes=5MB
environment=NODE_ENV="production"

; ---------------------------------------------------------------------------
; bot-presence
; ---------------------------------------------------------------------------
[program:bot-presence]
command=node dist/main.js
directory=/app/services/bot-presence
autostart=true
autorestart=true
startretries=10
priority=30
stdout_logfile=/home/container/logs/bot-presence.log
stdout_logfile_maxbytes=5MB
stderr_logfile=/home/container/logs/bot-presence-err.log
stderr_logfile_maxbytes=5MB
environment=NODE_ENV="production"

; ---------------------------------------------------------------------------
; schedule-emitter (backend-api)
; ---------------------------------------------------------------------------
[program:schedule-emitter]
command=node dist/scripts/schedule-emitter.js
directory=/app/services/backend-api
autostart=true
autorestart=true
startretries=10
priority=30
stdout_logfile=/home/container/logs/schedule-emitter.log
stdout_logfile_maxbytes=5MB
stderr_logfile=/home/container/logs/schedule-emitter-err.log
stderr_logfile_maxbytes=5MB
environment=NODE_ENV="production"

; ---------------------------------------------------------------------------
; monolith – backend-api web server (UI + REST API)
; ---------------------------------------------------------------------------
[program:monolith]
command=node dist/main.js
directory=/app/services/backend-api
autostart=true
autorestart=true
startretries=10
priority=40
stdout_logfile=/home/container/logs/monolith.log
stdout_logfile_maxbytes=5MB
stderr_logfile=/home/container/logs/monolith-err.log
stderr_logfile_maxbytes=5MB
environment=NODE_ENV="production"
